# openvpn-server-routing.service
[Unit]
Description=Enabling routing of OpenVPN traffic.
[Service]
ExecStart=/root/bin/vpn_route.sh
[Install]
WantedBy=multi-user.target
